# systemd/user/battery-watch.service
[Unit]
Description=Battery Watcher

[Service]
ExecStart=/home/pt/.local/bin/battery-watch.sh
Restart=always

[Install]
WantedBy=default.target
