# systemd/user/personal-hotplug.service
[Unit]
Description=display hotplug handler
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=oneshot
Environment=DISPLAY=:0
Environment=XAUTHORITY=%t/gdm/Xauthority
ExecStart=%h/.config/display/postswitch

