# systemd/user/time-brightness.service
[Unit]
Description=Apply time-based brightness profile
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/.config/display/time-brightness auto

